# src/debian/chronos-gui.service
[Unit]
Description=User interface for the Chronos Camera
After=chronos-video.service chronos-fpga.service
Requires=chronos-video.service chronos-fpga.service

[Service]
WorkingDirectory=/var/camera
EnvironmentFile=/etc/chronos-gui.conf
ExecStart=/usr/bin/camApp -qws
StandardOutput=journal+console
StandardError=inherit

[Install]
WantedBy=multi-user.target
